# tele_agent.service
[Unit]
Description=AI Agent Telegram Bot
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=%u
WorkingDirectory=%h/code/agent_mini
Environment="PATH=%h/code/agent_mini/venv/bin"
ExecStart=%h/code/agent_mini/venv/bin/python3 %h/code/agent_mini/tele_agent.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
